# yt-downloader-bot/rootfs/etc/s6-overlay/s6-rc.d/yt-bot/run
#!/command/with-contenv bashio

export BOT_TOKEN=$(bashio::config 'bot_token')
export LOCAL_API_URL="http://127.0.0.1:8081"

# Wait for telegram-bot-api to be ready
sleep 3

bashio::log.info "Starting YouTube Downloader Bot..."

cd /app
exec python3 -m bot.main
