# yt-downloader-bot/rootfs/etc/services.d/telegram-bot-api/run
#!/usr/bin/with-contenv bashio

exec telegram-bot-api \
    --api-id="${TELEGRAM_API_ID}" \
    --api-hash="${TELEGRAM_API_HASH}" \
    --local \
    --http-port=8081
